# app-crypt/oscam/files/oscam.initd
#!/sbin/openrc-run

depend() {
	before tvheadend
}

start() {
	ebegin "Starting OSCAM"
	start-stop-daemon --start --quiet --background \
		--make-pidfile --pidfile /var/run/oscam.pid \
		--exec /usr/bin/oscam -- ${OSCAM_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping OSCAM"
	start-stop-daemon --stop --quiet --pidfile /var/run/oscam.pid \
		--exec /usr/bin/oscam
	eend $?
}
